# systemd/e3dc-soh-measure.service
[Unit]
Description=E3DC SOH Messung

[Service]
Type=oneshot
WorkingDirectory=/opt/e3dcset-ui
EnvironmentFile=-%h/.config/e3dc-soh-monitor/env
ExecStart=/opt/e3dcset-ui/.venv/bin/python /opt/e3dcset-ui/measure.py --source timer
