# systemd/e3dc-soh-measure.service
[Unit]
Description=E3DC SOH Messung

[Service]
Type=oneshot
WorkingDirectory=/opt/e3dc-soh-monitor
EnvironmentFile=-%h/.config/e3dc-soh-monitor/env
ExecStart=/opt/e3dc-soh-monitor/.venv/bin/python /opt/e3dc-soh-monitor/measure.py --source timer
